# ansible/roles/power_management/templates/homelab-idle-monitor.service.j2
[Unit]
Description=HomeLab Idle Monitor
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/idle-monitor.sh
Restart=always
RestartSec=30

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target